# scripts/anyrouter-checkin.service
[Unit]
Description=AgentRouter Check-in (本机专用)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# 修改为你的项目绝对路径
WorkingDirectory=/home/decker/app/tools/anyrouter-check-in/anyrouter-check-in
ExecStart=/home/decker/app/tools/anyrouter-check-in/anyrouter-check-in/scripts/checkin-wrapper.sh
Environment=PATH=/home/decker/.local/bin:/usr/local/bin:/usr/bin:/bin

# 无界面环境，强制 headless 模式（timer 无 DISPLAY）
Environment=HEADLESS=1

[Install]
WantedBy=default.target
